# scripts/systemd/axon-compound-daemon.service
[Unit]
Description=Axon Agent Auto-Compound Daemon (自动复投奖励)
After=network-online.target axon-heartbeat-daemon.service
Wants=axon-heartbeat-daemon.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/axon-agent-scale
Environment=PYTHONPATH=/home/ubuntu/axon-agent-scale/current
ExecStart=/usr/bin/python3 /home/ubuntu/axon-agent-scale/scripts/compound.py \
    daemon \
    --state  /home/ubuntu/axon-agent-scale/state/deploy_state.json \
    --network /home/ubuntu/axon-agent-scale/current/configs/network.yaml \
    --agents  /home/ubuntu/axon-agent-scale/current/configs/agents.yaml \
    --config  /home/ubuntu/axon-agent-scale/current/configs/compound.yaml
Restart=on-failure
RestartSec=30
Environment="PYTHONUNBUFFERED=1"
StandardOutput=journal
StandardError=journal
SyslogIdentifier=axon-compound-daemon

[Install]
WantedBy=multi-user.target
